# data/init/systemd/default/process_master.service.erb
[Unit]
StopWhenUnneeded=true
Requires=<%= name %>.service
After=<%= name %>.service

[Service]
# this service is just a placeholder
ExecStart=/bin/bash -c 'while true; do sleep 10m; done'

[Install]
WantedBy=<%= name %>.service
