# lib/pkgr/data/distributions/debian/upstart/process.conf.erb
start on starting <%= name %>-<%= process_name %>
stop on stopping <%= name %>-<%= process_name %>
respawn

env PORT=PORT_NUM

exec <%= name %> run <%= process_name %> >> /var/log/<%= name %>/<%= process_name %>-PROCESS_NUM.log 2>&1
