# data/init/upstart/1.5/process_master.conf.erb
start on starting <%= name %>
stop on stopping <%= name %>
