# lib/pkgr/data/distributions/debian/sysv/process_master.erb
#!/bin/sh
# Init script for <%= name %>-<%= process_name %>
# Generated by pkgr.
# Implemented based on LSB Core 3.1:
#   * Sections: 20.2, 20.3
#
### BEGIN INIT INFO
# Provides:          <%= name %>-<%= process_name %>
# Required-Start:    <%= name %>
# Required-Stop:     <%= name %>
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       <%= process_name %> component of <%= name %>
### END INIT INFO

set -e

name="<%= name %>"
process_name="<%= process_name %>"

delegate() {
  for file in $(ls -1 /etc/init.d/${name}-${process_name}-* 2>/dev/null); do
    ${file} "${1}"
  done
}

case "$1" in
  start|stop|force-stop|status|restart) delegate "${1}" ;;
  *)
    echo "Usage: $SCRIPTNAME {start|stop|force-stop|status|restart}" >&2
    exit 3
  ;;
esac

exit $?